# etc/systemd/system/nm-unmetered-connection.service
[Unit]
Description=Check if the current NetworkManager connection is metered

[Service]
Type=oneshot
ExecStart=/sbin/nm-unmetered-connection.sh
